# launchd/com.imessage-reply-drafter-sticky.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Optional: always-on-top "sticky note" window showing the latest draft
  (sticky.py). To use it:
    1. Your Python needs a working Tk. The macOS command-line-tools Python
       ships a broken Tk 8.5 — install a good one, e.g.:
           brew install python-tk@3.13
       and point ProgramArguments at that interpreter (e.g. /opt/homebrew/bin/python3.13).
    2. Copy to ~/Library/LaunchAgents/com.imessage-reply-drafter-sticky.plist
       and replace the paths below with your own.
    3. launchctl load -w ~/Library/LaunchAgents/com.imessage-reply-drafter-sticky.plist
  The sticky needs NO Full Disk Access (it only reads latest_draft.json).
  The optional "paste into Messages" (⌘V) button needs Accessibility permission
  for the Python binary (System Settings → Privacy & Security → Accessibility).
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.imessage-reply-drafter-sticky</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/python3.13</string>
    <string>/Users/YOUR_USERNAME/imessage-reply-drafter/sticky.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/YOUR_USERNAME/imessage-reply-drafter</string>
  <key>EnvironmentVariables</key>
  <dict><key>PYTHONUNBUFFERED</key><string>1</string></dict>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>ThrottleInterval</key><integer>10</integer>
  <key>StandardOutPath</key><string>/Users/YOUR_USERNAME/imessage-reply-drafter/sticky.log</string>
  <key>StandardErrorPath</key><string>/Users/YOUR_USERNAME/imessage-reply-drafter/sticky.log</string>
</dict>
</plist>
